# etc/init/kiosk.conf
start on (filesystem and stopped udevtrigger)
stop on runlevel [06]

emits starting-x
respawn

exec sudo -u kiosk startx /etc/X11/Xsession /opt/kiosk.sh --